# config/redis.conf
[program:redis]
command=/usr/local/bin/redis-server
autostart=true
autorestart=true
user=django
redirect_stderr=true
stdout_logfile = /home/django/datein/logs/redis.log
startsecs=3